# nlw/service/st-userd/run
#!/usr/bin/env perl
# @COPYRIGHT@
use strict;
use warnings;
BEGIN { push @INC, "$ENV{ST_CURRENT}/nlw/lib" if $ENV{ST_CURRENT} }
use Socialtext::Daemontools;
Socialtext::Daemontools->RunSupervised(
    log_file => 'st-userd.log',
    port_name => 'userd_port',
    cb => sub {
        my $port = shift;
        my $exe = '/usr/bin/st-userd';
        if (!Socialtext::AppConfig->is_appliance()) {
            $exe = "$ENV{ST_CURRENT}/nlw/bin/st-userd";
        }
        exec $exe,'--port',$port
            or die "can't exec ${exe}: $!";
    }
);
